# mute_setting/www/cgi-bin/Update/updchk.timer
[Unit]
Description=timer updchk.service

[Timer]
# OnUnitActiveSec=10m
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
